# Default/0003-SyphaxOSGnome3/001_BuildPackagesScripts/0137-cups-2.1.4/cups.install
# arg 1:  the new package version
post_install() {
   useradd -c "SyphaxOS Print Service" -d /var/spool/cups -g lp -s /bin/false -u 9 lp
   groupadd -g 19 lpadmin
   usermod -a -G lpadmin root
   rm -rf /tmp/cupsinit
   gtk-update-icon-cache
   systemctl enable org.cups.cupsd
}
